# homeassistant-addon/rootfs/etc/s6-overlay/s6-rc.d/dispatcharr/run
#!/usr/bin/with-contenv bashio
bashio::log.info "RUN SCRIPT IS EXECUTING - THIS IS A TEST"
sleep 10
bashio::log.info "TEST COMPLETE"
